# docker/supervisord.swoole.conf
[program:swoole]
command=php artisan octane:start --server=swoole --host=0.0.0.0 --port=8000
user=%(ENV_USER)s
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/swoole.log
stopwaitsecs=3600
